# deploy/lila-rs-manager.service
[Unit]
# The Rust little-living-apps manager, as a CUTOVER replacement for the TS `lila-manager.service` on
# the live host (not side-by-side — that is deploy/lila-rs@.service). One self-contained binary at
# /opt/lila/bin/lila; env from /etc/lila/lila-rs.env. Restart=always makes a crash or a `/backend`
# swap (which exits 0) come straight back. Per the rewrite's fresh-schema design the snapshot/session
# do NOT carry over from the TS manager — memory dir is a clean slate (operator's choice).
Description=little-living-apps (Rust) manager — cutover
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/etc/lila/lila-rs.env
ExecStart=/opt/lila/bin/lila run
Restart=always
RestartSec=2
Environment=RUST_BACKTRACE=1
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
