# deploy/lila-manager.service
# systemd unit for the little-living-apps manager. bootstrap.sh substitutes the __PLACEHOLDERS__
# and installs this to /etc/systemd/system/lila-manager.service. Keeps the manager alive across
# crashes and reboots; a cold restart is lossless (snapshot + persistent memory).
#
# The manager opens NO inbound port — it reaches Telegram by outbound long-poll — so there is
# nothing here to expose. The disposable VM is the security boundary.
[Unit]
Description=little-living-apps manager (Telegram -> Codex manager + workers)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=/etc/lila/lila.env
# Worker self-validation: playwright lives at a fixed location (bootstrap §4b), exposed here so that
# Codex workers — child processes of this manager, inheriting its env — can `require("playwright")`
# in interactive validation scripts (and `npx playwright …`) from any cwd. Harmless to the manager's
# own node process: local node_modules still take precedence; this is only a fallback resolution path.
Environment=NODE_PATH=/opt/lila/tooling/node_modules
# mise exec resolves the pinned Ruby/Node (and the global codex bin) from .mise.toml in WorkingDirectory.
ExecStart=__MISE__ exec -- npm start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
