# deploy/lila-manager@.service
# systemd TEMPLATE unit for a little-living-apps manager INSTANCE (Rust `lila` binary). One instance
# = one brain = one app (the model never changes: each manager owns exactly one app). This is how you
# run several living apps on a single host — each a fully independent stack, isolated by its own
# per-instance env file.
#
# bootstrap.sh (or bin/new-instance) installs this template ONCE per host, substituting the
# host-level __PLACEHOLDERS__. The %i instance name then selects everything per-instance via
# /etc/lila/%i.env (its own Telegram bot token, WORKSPACE_DIR, CODEX_HOME, ports, domain). Enable an
# instance with:
#
#     sudo systemctl enable --now lila-manager@<name>
#
# The manager opens NO inbound port — it reaches Telegram by outbound long-poll. The disposable VM is
# the security boundary. A cold restart is lossless (crash snapshot + persistent memory survive).
[Unit]
Description=little-living-apps manager [%i] (Telegram -> manager + ephemeral workers)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=/var/lib/lila
EnvironmentFile=/etc/lila/%i.env
# Tell the manager which app unit it owns, so the prompt hands workers the RIGHT restart command.
Environment=LILA_APP_SERVICE=lila-app@%i
# The instance name itself, inherited by workers — lila-new-app reads it to (re)scaffold the RIGHT
# app (/srv/%i under lila-app@%i), never a sibling instance's.
Environment=LILA_INSTANCE=%i
# Worker self-validation: Playwright lives at a fixed, node-version-independent location (bootstrap
# §4b), exposed here so workers — children of this manager, inheriting its env — can
# require("playwright") and `npx playwright …` from any cwd.
Environment=NODE_PATH=/opt/lila/tooling/node_modules
# The codex/claude CLIs are mise-managed node scripts; the native binary gets no mise shell hook, so
# put the mise node bin dir (holds BOTH `codex`/`claude` AND the `node` they shebang to) on PATH.
# bootstrap.sh substitutes __NODEBIN__ from `mise which codex`.
Environment=PATH=__NODEBIN__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# The self-contained release binary, built + installed to /opt/lila/bin by bootstrap.sh.
ExecStart=/opt/lila/bin/lila run
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
